# roles/client/files/viewloop.service
[Unit]
Description="Starts VMware loop.sh script on login"

[Service]
ExecStart=/bin/sh /loop.sh

[Install]
WantedBy=multi-user.target